# hci-proxy.service
[Unit]
Description=Bluetooth HCI-over-Serial Proxy
After=bluetooth.target
Conflicts=bluetooth.service
# The proxy exits when the adapter re-enumerates (USB reset) so it can re-resolve the
# new index; recovery is a restart loop. Disable the start-rate limit so an
# enumeration storm can't wedge the unit in "start-limit-hit".
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStartPre=-/usr/bin/systemctl stop bluetooth.service
ExecStart=/usr/bin/python3 /usr/local/bin/hci-proxy.py
Restart=on-failure
RestartSec=5

# Needs raw Bluetooth socket access and connection to QEMU's virtio-serial socket
# Run as root for simplicity (HCI_CHANNEL_USER requires CAP_NET_ADMIN + CAP_NET_RAW)

[Install]
WantedBy=multi-user.target
